# config/rc/init.d/gpg-agent
#!/sbin/openrc-run
# GnuPG agent (provides gpg + ssh auth sockets under $XDG_RUNTIME_DIR/gnupg).
#
# NOT supervise-daemon: gpg-agent 2.5.20 has no foreground+listen mode a
# supervisor can track. --daemon always forks (orphans to init), --server is
# stdin-only (no sockets), and --deprecated-supervised needs socket-activation
# FDs OpenRC doesn't hand out. gpg-agent self-relaunches on demand anyway, so we
# just launch it GnuPG's own way (with the session bus set) and kill it on stop.

# ssh signatures can't carry a display/bus per request the way gpg does (the
# ssh-agent protocol has no OPTION to pass tty/display), so the agent must inherit
# the session bus at launch or pinentry-gnome3 can't reach the GNOME/gcr prompter
# -- ssh ops then fail with "agent refused operation", worst for the confirm/
# passphrase keys in ~/.gnupg/sshcontrol. On this Wayland GNOME box every session
# daemon reaches the GUI purely over the session bus (no DISPLAY/WAYLAND), so this
# one var is all pinentry needs -- same idiom as /etc/user/init.d/gsd and
# xdg-desktop-portal. Do NOT import the whole gnome-session environ instead: it
# carries a PATH without /lib/rc/bin, which breaks OpenRC's own ebegin/eend.
DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
export DBUS_SESSION_BUS_ADDRESS

description="GnuPG agent (gpg + ssh auth sockets)"

depend() {
	after dbus
}

start() {
	ebegin "Launching gpg-agent"
	# Be the sole launcher: evict any agent auto-started on demand (which may
	# lack the session bus above), then launch it explicitly with gpgconf so the
	# service -- not a stray connection probe -- owns the invocation. (Using
	# `gpg-connect-agent /bye` instead only autostarts the agent as a side effect
	# of connecting, with an 8s connect-wait.)
	gpgconf --kill gpg-agent 2>/dev/null
	gpgconf --launch gpg-agent
	eend $? "Failed to launch gpg-agent"
}

stop() {
	ebegin "Stopping gpg-agent"
	gpgconf --kill gpg-agent
	eend $? "Failed to stop gpg-agent"
}

status() {
	# --no-autostart: probe without launching, so a status check can't start it.
	# NB: gpg-connect-agent exits 0 even when no agent is running, so key off the
	# 'OK' response, not the exit code -- otherwise status always reports running.
	if gpg-connect-agent --no-autostart 'GETINFO version' /bye 2>/dev/null | grep -q '^OK'; then
		einfo "gpg-agent is running"
		return 0
	fi
	einfo "gpg-agent is not running"
	return 3
}
